# resources/etc/systemd/system/instruqt-configure-aws.path
[Unit] 
Description = instruqt environment settings
Wants = instruqt-configure-aws.service
[Path]
PathChanged = /etc/profile.d/instruqt-env.sh
[Install]
WantedBy = multi-user.target
