# resources/etc/systemd/system/instruqt-configure-gcloud.service
[Unit]
Description=Instruqt GCloud credentials configuration

[Service]
ExecStart=/usr/local/bin/instruqt-configure-gcloud
Type=simple
