# resources/etc/systemd/system/instruqt-configure-aws.service
[Unit]
Description=Instruqt AWS credentials configuration

[Service]
ExecStart=/usr/local/bin/instruqt-configure-aws
Type=simple

